# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1024x768x24
autorestart=true
user=firefox

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=:1
autorestart=true
user=firefox

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -listen localhost -xkb -forever
autorestart=true
user=firefox

[program:novnc]
command=/usr/bin/python3 /opt/websockify/websockify.py --web /opt/novnc 3000 localhost:5901
autorestart=true
user=firefox

[program:firefox]
command=/usr/bin/firefox --display=:1 --no-sandbox --enable-features=VaapiVideoDecoder --use-gl=egl --enable-gpu-rasterization https://whatismyipaddress.com/
environment=DISPLAY=:1,HOME=/home/firefox
autorestart=true
user=firefox
startsecs=10